# admin/supervisord.conf
[supervisord]
logfile = /tmp/supervisord.log

[program:PS4M]
autorestart=true
command=python -u main.py
directory=/home/ubuntu/PS4M
redirect_stderr=true
stdout_logfile=/home/ubuntu/logs/main.txt
